# .mutt/crontab
*/5 * * * * /Users/tommy/.mutt/offlineimapsync.sh
